# systemd/user/rdm.service
[Unit]
Description=RTags daemon

Requires=rdm.socketp

[Service]
Type=simple
ExecStart= -v --inactivity-timeout 300
